# newland-site/systemd_newland.service
[Unit]
Description=Newland site
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/newland-site.pid
Environment=CATALINA_PID=/var/run/newland-site.pid
Environment=TOMCAT_JAVA_HOME=/usr/lib/jvm/java-7-openjdk
Environment=CATALINA_HOME=/usr/share/tomcat7
Environment=CATALINA_BASE=/usr/share/newland-site
Environment=CATALINA_OPTS=

ExecStart=/usr/bin/jsvc \
            -Dcatalina.home=${CATALINA_HOME} \
            -Dcatalina.base=${CATALINA_BASE} \
            -Djava.io.tmpdir=/var/tmp/newland-site/temp \
            -cp /usr/share/java/commons-daemon.jar:/usr/share/java/eclipse-ecj.jar:${CATALINA_HOME}/bin/bootstrap.jar:${CATALINA_HOME}/bin/tomcat-juli.jar \
            -user tomcat7 \
            -java-home ${TOMCAT_JAVA_HOME} \
            -pidfile /var/run/newland-site.pid \
            -errfile /var/log/newland-site/catalina.err \
            -outfile /var/log/newland-site/catalina.out \
            $CATALINA_OPTS \
            org.apache.catalina.startup.Bootstrap

ExecStop=/usr/bin/jsvc \
            -pidfile /var/run/newland-site.pid \
            -stop \
            org.apache.catalina.startup.Bootstrap

[Install]
WantedBy=multi-user.target

